# root/etc/services.d/rutorrent/run
#!/usr/bin/with-contenv bash

_term() {
	echo "Caught SIGTERM for rtorrent!"
	killall -s SIGINT rtorrent 2>/dev/null
	sleep 1
}

trap _term SIGTERM

dtach -N /detach_sess/.rtorrent \
	s6-setuidgid abc /usr/bin/rtorrent \
		-n -o import=/config/rtorrent/rtorrent.rc &

wait

sleep 1s
